# systemd/autograder-worker-bulk.service
[Unit]
Description=Autograder Celery Worker Bulk (whatsapp_bulk, concurrency=1)
After=postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
Type=exec
User=autograder
Group=autograder
WorkingDirectory=/opt/autograder/autograder-back
EnvironmentFile=/opt/autograder/.env
ExecStart=/opt/autograder/autograder-back/.venv/bin/celery -A app.celery_app worker --loglevel=info --concurrency=1 -Q whatsapp_bulk
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autograder-worker-bulk

[Install]
WantedBy=multi-user.target
