# systemd/autograder-discord.service
[Unit]
Description=Autograder Discord Bot
After=postgresql.service redis-server.service
Requires=postgresql.service redis-server.service

[Service]
Type=exec
User=autograder
Group=autograder
WorkingDirectory=/opt/autograder/autograder-back
EnvironmentFile=/opt/autograder/.env
ExecStart=/opt/autograder/autograder-back/.venv/bin/python -m app.discord_bot
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=autograder-discord

[Install]
WantedBy=multi-user.target
